# database/crontab
*/1 * * * * cd /app && /usr/local/bin/python -m core.main.py >> cron.log 2>&1
*/60 * * * * /bin/cat /dev/null > /app/cron.log
